# mistralvibe/rootfs/etc/services.d/vibe-terminal/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Stop the container as soon as the terminal service exits, so HA Supervisor
# can restart us cleanly instead of leaving an init shell running.
# ==============================================================================
bashio::log.warning "vibe-terminal exited with code $1 – stopping add-on"
/run/s6/basedir/bin/halt
